# uwsgiConfig/uwsgi.service
[Unit]
description=Service file for Uwsgi application
[Service]
type=notify
ExecStart=/usr/local/bin/uwsgi_run.sh
Restart=on-abnormal

